# script_10.service
[Unit]
Description = New file in watch
After = network.target
[Service]
ExecStart = /usr/sbin/script_10.sh
[Install]
WantedBy = multi-user.target
